# deployment/systemd/mule-forward.service
[Unit]
Description=Mule encrypted TCP forwarder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mule
Group=mule
ExecStart=/usr/local/bin/mule forward \
  --peer host-a.example.org:4400 \
  --forward-id host-b \
  --listen web=127.0.0.1:3100 \
  --listen ssh=127.0.0.1:3101 \
  --secret-file /etc/mule/b-to-a.key \
  --log-format json
Restart=on-failure
RestartSec=3

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/mule
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
