# deployment/systemd/mule-forward.service
[Unit]
Description=Mule encrypted TCP tunnel agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mule
Group=mule
ExecStart=/usr/local/bin/mule agent \
  --server host-a.example.org:4400 \
  --agent-id host-b \
  --secret-file /etc/mule/host-b.key \
  --forward web=127.0.0.1:3100 \
  --reverse ssh=127.0.0.1:22 \
  --log-format json
Restart=on-failure
RestartSec=3

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/mule
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
